# ops/systemd/ollama.service.d/10-keepalive.conf
[Service]
Environment=OLLAMA_NUM_PARALLEL=1
Environment=OLLAMA_KEEP_ALIVE=15m
